# local/cron.d/99-enable-scheduler-in-typo3
*/5  *  * * *  TYPO3_CONTEXT=Development /var/www/html/web/typo3/cli_dispatch.phpsh scheduler > /dev/null 2>&1
